# configs/supervisor.d/cron.ini
[program:cron]
command=/usr/bin/crond -n
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cron.log
